# network/systemd/archcfg-reset-resolved-if-stub.service
[Unit]
Description=Reset systemd-resolved after a VPN restores the stub resolver
After=systemd-resolved.service

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/archcfg-reset-resolved-if-stub
