# live/airootfs/etc/systemd/system/greetd.service.d/10-vtrgb.conf
[Service]
# Theme-only step; never block login if console palette setup fails.
ExecStartPre=-/usr/bin/setvtrgb /etc/vtrgb
