# live/airootfs/etc/systemd/system/greetd.service.d/10-vtrgb.conf
[Service]
# Give the kernel console a moment to settle at its final KMS size before tuigreet draws.
ExecStartPre=/usr/bin/sleep 1
# Theme-only step; never block login if console palette setup fails.
ExecStartPre=-/usr/bin/setvtrgb /etc/vtrgb
