# wayland/systemd/user/hyprsunset-day.service
[Unit]
Description=Reset Hyprsunset during daytime

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprsunset-set identity
