# wayland/systemd/user/hyprsunset-night.service
[Unit]
Description=Set Hyprsunset strong warm temperature for night

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprsunset-set temperature 3500
