# wayland/systemd/user/xdg-desktop-portal-hyprland.service.d/share-picker.conf
[Service]
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin"
